# debian/rawdns.upstart
description "raw DNS interface to the Docker API"
author "Tianon Gravi <admwiggin@gmail.com>"

start on local-filesystems
stop on runlevel [!2345]

respawn

script
	BINARY=/usr/bin/$UPSTART_JOB
	CONFIG=/etc/$UPSTART_JOB.json
	if [ -f /etc/default/$UPSTART_JOB ]; then
		. /etc/default/$UPSTART_JOB
	fi
	exec "$BINARY" "$CONFIG"
end script

# vim:set ft=upstart:
